# system_files/etc/profile.d/kawos-default-apps.sh
#!/usr/bin/env bash

# Distro defaults only: preserve any value explicitly selected by the user.
export TERMINAL="${TERMINAL:-ghostty}"
export BROWSER="${BROWSER:-helium}"
export EDITOR="${EDITOR:-nvim}"
export VISUAL="${VISUAL:-nvim}"
